# templates/default/supervisor.conf.erb
[program:<%= @app_name %>]
command=<%= @env_home %>/bin/gunicorn_django -c /etc/gunicorn/<%= @app_name %>.py
environment=HOME='<%= @env_home %>'
autostart=true
autorestart=true
directory=<%= @container_dir %>
